# scripts/com.claude-notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-notifier</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ab180-minchae-mbp/Applications/ClaudeNotifier.app/Contents/MacOS/ClaudeNotifier</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardOutPath</key>
    <string>/Users/ab180-minchae-mbp/.claude-notifier/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ab180-minchae-mbp/.claude-notifier/stderr.log</string>
</dict>
</plist>
